# s3-backup/rootfs/etc/services.d/amazon-s3-backup/finish
#!/usr/bin/with-contenv bashio
# Cleanup logic (if any)
bashio::log.info "Finished Amazon S3 Backup."