# s3-backup/rootfs/etc/services.d/amazon-s3-backup/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Amazon S3 Backup
# ==============================================================================
bashio::log.info "Finishing Amazon S3 Backup..."
s6-svscanctl -t /var/run/service
bashio::log.info "Finishing Amazon S3 Backup... Done!"

